# mopidy-scan.timer
[Unit]
Description=Scans local library nightly

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
